# discord-bot/deploy/characterai-bot.service
[Unit]
Description=Neural Cord Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/neuralcord/discord-bot
ExecStart=/usr/bin/python3 src/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neuralcord
Environment="PYTHONUNBUFFERED=1"

# Uncomment and modify if you need to set environment variables
# EnvironmentFile=/opt/neuralcord/discord-bot/.env

[Install]
WantedBy=multi-user.target 